# autossh-simcom.service.sample
############################################################################################
# systemd unit sample: autossh-simcom.service
############################################################################################
# Reverse-tunnel service gated on modem-data.service. Starts only after the SIM8262E
# QMI session has an IP on wwan0. This is the SIMCom replacement for the legacy
# maintenance.sh + wvdial + manual autossh invocation used on PLS8-E deployments.
#
# The jumphost pattern mirrors the mobileid013/005 deployments: remote port is
# 19<hostnumber>, forwards to localhost:5000 on the Pi.
#
# Deployment:
#   sudo cp autossh-simcom.service.sample /etc/systemd/system/autossh-simcom.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now autossh-simcom.service
#
# Adjust JUMPHOST, JUMP_USER, JUMP_PORT, REMOTE_PORT, and LOCAL_PORT below to match
# the deployment environment.
############################################################################################

[Unit]
Description=AutoSSH reverse tunnel over SIMCom wwan0
After=modem-data.service network-online.target
Requires=modem-data.service

[Service]
Type=simple
Environment=JUMPHOST=ssh.mobileid.host
Environment=JUMP_USER=mid
Environment=JUMP_PORT=5000
Environment=REMOTE_PORT=19000
Environment=LOCAL_PORT=5000
Environment=SSH_IDENTITY=/root/.ssh/id_rsa
Environment=AUTOSSH_PIDFILE=/run/autossh-simcom.pid
Environment=AUTOSSH_GATETIME=30
Environment=AUTOSSH_POLL=60

ExecStart=/usr/bin/autossh -M 0 \
    -o "ServerAliveInterval 30" \
    -o "ServerAliveCountMax 3" \
    -o "ExitOnForwardFailure=yes" \
    -o "StrictHostKeyChecking=yes" \
    -o "UserKnownHostsFile=/root/.ssh/known_hosts" \
    -N \
    -R ${REMOTE_PORT}:localhost:${LOCAL_PORT} \
    -i ${SSH_IDENTITY} \
    -p ${JUMP_PORT} \
    ${JUMP_USER}@${JUMPHOST}

Restart=on-failure
RestartSec=60

[Install]
WantedBy=multi-user.target
